# launchd/com.zhb.abell.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.zhb.abell</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/zhb/Documents/a-bell/.venv/bin/python</string>
    <string>-m</string><string>abell.main</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/zhb/Documents/a-bell</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/zhb/Documents/a-bell/data/abell.log</string>
  <key>StandardErrorPath</key><string>/Users/zhb/Documents/a-bell/data/abell.err.log</string>
</dict>
</plist>
